# crontab
* * * * * /data/content-recommendation update >> /var/log/cron.log 2>&1
0 0 * * * python3 /data/dbutils/utils.py --cron >> /var/log/cron.log 2>&1
