# crontab
*/5 * * * * ./content-recommendation update
0 0 * * * python3 utils.py